# contrib/yzu-wlan.service
[Unit]
Description=YZU Campus Network Auto-Login
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/yzu-wlan

[Install]
WantedBy=default.target
